# rpm/udhcpd.service
[Unit]
Description=udhcpcd DHCP server
[Service]
ExecStart=/usr/sbin/udhcpd -f
